# deploy/services/taskq-dispatcher.timer
[Unit]
Description=Run task queue dispatcher every 2 minutes

[Timer]
OnCalendar=*:0/2
Persistent=true

[Install]
WantedBy=timers.target
